# supervisord.conf
; XPS Intelligence Platform — All-in-one supervisord process manager
; Manages both the FastAPI backend and the Next.js frontend within a single
; Railway container.
;
; Port layout:
;   $PORT  (Railway-assigned, typically 3000) → Next.js  (public)
;   8000                                       → FastAPI  (localhost only)
;
; The Next.js server proxies /api/v1/* → http://127.0.0.1:8000/api/v1/*
; via the rewrites() in next.config.js (BACKEND_INTERNAL_URL env var).

[supervisord]
nodaemon=true
logfile=/dev/stdout
logfile_maxbytes=0
loglevel=info

; ── FastAPI backend (uvicorn) ─────────────────────────────────────────────────
[program:backend]
command=uvicorn apps.backend.main:app
    --host 127.0.0.1
    --port 8000
    --workers 2
    --log-level info
directory=/app
environment=PYTHONPATH=/app
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Next.js frontend (standalone) ────────────────────────────────────────────
[program:frontend]
command=node apps/frontend/server.js
directory=/app
environment=
    PORT=%(ENV_PORT)s,
    HOSTNAME=0.0.0.0,
    NODE_ENV=production,
    BACKEND_INTERNAL_URL=http://127.0.0.1:8000,
    NEXT_TELEMETRY_DISABLED=1
autostart=true
autorestart=true
startretries=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
